# deploy/smartgrow.service
# /home/smartgrow/deploy/smartgrow.service
# Systemd unit файл
#
# Встановлення:
#   sudo cp deploy/smartgrow.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable smartgrow
#   sudo systemctl start smartgrow

[Unit]
Description=SmartGrow Mini — Autonomous Greenhouse Controller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/smartgrow
ExecStart=/home/smartgrow/venv/bin/python /home/smartgrow/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
